# com.nova.assistant.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" 
"http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nova.assistant</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/sohamdhande/Docs_Local/NOVA/venv/bin/python3</string>
        <string>/Users/sohamdhande/Docs_Local/NOVA/main.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/sohamdhande/Docs_Local/NOVA</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/sohamdhande/Docs_Local/NOVA/logs/nova.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sohamdhande/Docs_Local/NOVA/logs/nova_error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
